# crontab
*/15 * * * * nextcloud  /usr/bin/php -f /srv/www/cron.php
