# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -c /etc/nginx/nginx.conf

[program:fpm]
command=php-fpm7.2 -F

[program:crond]
command=cron -f
